# debian/sensord.conffiles
/etc/default/sensord
/etc/init.d/sensord
/etc/logcheck/ignore.d.server/sensord
/etc/logcheck/ignore.d.workstation/sensord
